# deploy/alibaba/describeops.service
[Unit]
Description=DescribeOps Django backend
After=network.target

[Service]
Type=simple
User=describeops
Group=describeops
WorkingDirectory=/opt/describeops/backend
EnvironmentFile=/opt/describeops/backend/.env
ExecStart=/opt/describeops/backend/.venv/bin/gunicorn describeops_backend.wsgi:application --bind 127.0.0.1:8000 --workers 2 --threads 4 --timeout 300 --access-logfile - --error-logfile -
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
